# dot_config/systemd/user/rclone-onedrive.timer
[Unit]
Description=Backup to OneDrive using rclone

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
