# dot_config/systemd/user/touchpad_toggle.service
[Unit] 
Description=Touchpad Toggle

[Service] 
ExecStart=/home/ye/bin/touchpad_toggle.sh

[Install] 
WantedBy=default.target
