# dot_config/systemd/user/timers.target.wants/backup-dotfiles.timer
[Unit] 
Description=Run backup dotfiles daily and on boot 

[Timer] 
OnBootSec=10min 
OnUnitActiveSec=1d

[Install] 
WantedBy=timers.target
