# dot_config/systemd/user/timers.target.wants/backup-enabled-systemd-units.timer
[Unit]
Description=Run backup enabled systemd units list daily and on boot

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
